# mbot/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- from "mbot/map.jinja" import server with context %}

[program:mbot_{{ service.name }}]
command=/srv/mbot/bin/mbot --config-path={{ server.dir.config }}
stdout_logfile=/var/log/mbot/mbot.log
stderr_logfile=/var/log/mbot/mbot_error.log
user=mbot
autostart=true
autorestart=true
